# rootfs/etc/services.d/ca350/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

CONFIG_PATH=/data/options.json

SOCAT="$(bashio::config 'REMOTE_DEVICE' 'false')"
SERIAL_PORT="$(bashio::config 'SERIAL_PORT')"

# Start SOCAT process
if [[ "${SOCAT,,}" == "true" ]]; then
    REMOTE_IP="$(bashio::config 'REMOTE_IP' 'localhost')"
    REMOTE_PORT="$(bashio::config 'REMOTE_PORT' '502')"

    bashio::log.info "create serial device over ethernet with socat for remote device $REMOTE_IP:$REMOTE_PORT"
    SERIAL_PORT=/dev/comfoair
    socat -d -d pty,link="$SERIAL_PORT",raw,group-late=dialout,mode=660 tcp:"$REMOTE_IP":"$REMOTE_PORT" &
else
    bashio::log.info "Using local serial device '$SERIAL_PORT'"
fi

bashio::log.info "Creating config.ini"

RS485_PROTOCOL="$(bashio::config 'RS485_PROTOCOL' 'false')"
REFRESH_INTERVAL="$(bashio::config 'REFRESH_INTERVAL')"
ENABLE_PC_MODE="$(bashio::config 'ENABLE_PC_MODE')"
DEBUG="$(bashio::config 'DEBUG' 'false')"

MQTT_SERVER="$(bashio::config 'MQTT_SERVER' 'localhost')"
MQTT_PORT="$(bashio::config 'MQTT_PORT' '1883')"
MQTT_KEEPALIVE="$(bashio::config 'MQTT_KEEPALIVE' '45')"
MQTT_USER="$(bashio::config 'MQTT_USER')"
MQTT_PASSWORD="$(bashio::config 'MQTT_PASSWORD')"

HA_ENABLE_AUTO_DISCOVERY_SENSORS="$(bashio::config 'HA_ENABLE_AUTO_DISCOVERY_SENSORS' 'true')"
HA_ENABLE_AUTO_DISCOVERY_CLIMATE="$(bashio::config 'HA_ENABLE_AUTO_DISCOVERY_CLIMATE' 'true')"
HA_AUTO_DISCOVERY_DEVICE_ID="$(bashio::config 'HA_AUTO_DISCOVERY_DEVICE_ID' 'ca350')"
HA_AUTO_DISCOVERY_DEVICE_NAME="$(bashio::config 'HA_AUTO_DISCOVERY_DEVICE_NAME' 'CA350')"
HA_AUTO_DISCOVERY_DEVICE_MANUFACTURER="$(bashio::config 'HA_AUTO_DISCOVERY_DEVICE_MANUFACTURER' 'Zehnder')"
HA_AUTO_DISCOVERY_DEVICE_MODEL="$(bashio::config 'HA_AUTO_DISCOVERY_DEVICE_MODEL' 'ComfoAir 350')"

DEVICE_FANOUT_ABSENT="$(bashio::config 'DEVICE_FANOUT_ABSENT' '15')"
DEVICE_FANOUT_LOW="$(bashio::config 'DEVICE_FANOUT_LOW' '35')"
DEVICE_FANOUT_MID="$(bashio::config 'DEVICE_FANOUT_MID' '50')"
DEVICE_FANOUT_HIGH="$(bashio::config 'DEVICE_FANOUT_HIGH' '70')"
DEVICE_FANIN_ABSENT="$(bashio::config 'DEVICE_FANIN_ABSENT' '20')"
DEVICE_FANIN_LOW="$(bashio::config 'DEVICE_FANIN_LOW' '40')"
DEVICE_FANIN_MID="$(bashio::config 'DEVICE_FANIN_MID' '55')"
DEVICE_FANIN_HIGH="$(bashio::config 'DEVICE_FANIN_HIGH' '75')"
DEVICE_SET_FAN_LEVEL_AT_START="$(bashio::config 'DEVICE_SET_FAN_LEVEL_AT_START' 'true')"

# build config file
sed \
    -e "s|SerialPort=.*|SerialPort=${SERIAL_PORT}|"  \
    -e "s|RS485_protocol=.*|RS485_protocol=${RS485_PROTOCOL}|"  \
    -e "s|refresh_interval=.*|refresh_interval=${REFRESH_INTERVAL}|"  \
    -e "s|enablePcMode=.*|enablePcMode=${ENABLE_PC_MODE}|"  \
    -e "s|debug=.*|debug=${DEBUG}|"  \
    -e "s|MQTTServer=.*|MQTTServer=${MQTT_SERVER}|"  \
    -e "s|MQTTPort=.*|MQTTPort=${MQTT_PORT}|"  \
    -e "s|MQTTKeepalive=.*|MQTTKeepalive=${MQTT_KEEPALIVE}|"  \
    -e "s|MQTTUser=.*|MQTTUser=${MQTT_USER}|"  \
    -e "s|MQTTPassword=.*|MQTTPassword=${MQTT_PASSWORD}|"  \
    -e "s|HAEnableAutoDiscoverySensors=.*|HAEnableAutoDiscoverySensors=${HA_ENABLE_AUTO_DISCOVERY_SENSORS}|"  \
    -e "s|HAEnableAutoDiscoveryClimate=.*|HAEnableAutoDiscoveryClimate=${HA_ENABLE_AUTO_DISCOVERY_CLIMATE}|"  \
    -e "s|HAAutoDiscoveryDeviceId=.*|HAAutoDiscoveryDeviceId=${HA_AUTO_DISCOVERY_DEVICE_ID}|"  \
    -e "s|HAAutoDiscoveryDeviceName=.*|HAAutoDiscoveryDeviceName=${HA_AUTO_DISCOVERY_DEVICE_NAME}|"  \
    -e "s|HAAutoDiscoveryDeviceManufacturer=.*|HAAutoDiscoveryDeviceManufacturer=${HA_AUTO_DISCOVERY_DEVICE_MANUFACTURER}|"  \
    -e "s|HAAutoDiscoveryDeviceModel=.*|HAAutoDiscoveryDeviceModel=${HA_AUTO_DISCOVERY_DEVICE_MODEL}|"  \
    -e "s|FanOutAbsent=.*|FanOutAbsent=${DEVICE_FANOUT_ABSENT}|" \
    -e "s|FanOutLow=.*|FanOutLow=${DEVICE_FANOUT_LOW}|" \
    -e "s|FanOutMid=.*|FanOutMid=${DEVICE_FANOUT_MID}|" \
    -e "s|FanOutHigh=.*|FanOutHigh=${DEVICE_FANOUT_HIGH}|" \
    -e "s|FanInAbsent=.*|FanInAbsent=${DEVICE_FANIN_ABSENT}|" \
    -e "s|FanInLow=.*|FanInLow=${DEVICE_FANIN_LOW}|" \
    -e "s|FanInMid=.*|FanInMid=${DEVICE_FANIN_MID}|" \
    -e "s|FanInHigh=.*|FanInHigh=${DEVICE_FANIN_HIGH}|" \
    -e "s|SetUpFanLevelsAtStart=.*|SetUpFanLevelsAtStart=${DEVICE_SET_FAN_LEVEL_AT_START}|" \
    /opt/hacomfoairmqtt/config.ini.dist >  /opt/hacomfoairmqtt/config.ini

bashio::log.info "Starting ca350.py"

# Start processing the ComfoAir
python /opt/hacomfoairmqtt/ca350.py